# data_processing/archive/data-processing.service
[Unit]
Description=Flow and plot new data summary files

[Service]
Type=oneshot
# Pass the specific file path that triggered the event
ExecStart=/full/path/to/data_processing.py "$SYSTEMD_PATH_FILENAME"
WorkingDirectory=/full/path/to/data_processing
User=your_username
Group=your_group
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target