# myproject/myproject_server.service
[Unit]
Description=MyProject Server
# Wait for network before starting
After=network.target
# Stop retrying after 3 crashes in 60 seconds
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=simple
# Run as a dedicated limited user, not root
User=myproject
WorkingDirectory=/etc/myproject
# Load secrets separately, not from git (temp solution define an env variable here)
#EnvironmentFile=/etc/myproject/secrets.env
Environment="MYPROJECT_API_KEY=your_secret_key"
# Tell config.py exactly where the config file is
Environment=MYPROJECT_CONFIG=/etc/myproject/config.yaml
# Use the venv Python, not the system Python
ExecStart=/opt/myproject/venv/bin/python -m myproject.server
# Only restart on crashes, not on clean stops
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
