# data_processing/data-processing@.service
[Unit]
Description=Process summary file %I
# Prevent the same file from being queued twice simultaneously
# (e.g. if close_write and moved_to both fire for the same file)

[Service]
Type=oneshot
ExecStart=python3 /path/to/data_processing.py "%I"
StandardOutput=journal
StandardError=journal